# profile/ihorwrt/files/etc/uci-defaults/99_custom
P=3X4a8A4v echo -e "$P\n$P" | passwd
uci -q batch << EOI
#System - System - General Settings
	#Hostname
	set system.@system[0].hostname='__NAME__'
#Network - DHCP and DNS - General Settings
	#General Settings 
        	#Local server
        	set dhcp.@dnsmasq[0].local='/ihor.local/'
        	#Local domain
        	set dhcp.@dnsmasq[0].domain='ihor.local'
#Network - Interfaces - Devices
        #br-lan
        delete network.@device[0]
        #eth1
        delete network.@device[-1]
        #VPN - Configure - Enable IPv6
	add network device
	set network.@device[-1].name='VPN'
	set network.@device[-1].ipv6='0'
#Network - Interfaces - Interfaces
	#LAN - Edit - General Settings
        	delete network.lan
	#WAN - Edit - General Settings
		set network.wan=interface
		set network.wan.device='eth0'
		set network.wan.proto='static'
		set network.wan.ipaddr='__WRTLAB_WAN_IP__'
		set network.wan.netmask='__WRTLAB_WAN_MASK__'
		set network.wan.gateway='__WRTLAB_WAN_GW__'
	#WAN - Edit - Advanced Settings
		#Use custom DNS servers
		set network.wan.dns='__WRTLAB_WAN_DNS__'
	#VPN - Edit - Advanced Settings
		set network.VPN=interface
		set network.VPN.proto='wireguard'
		add network wireguard_VPN
        	set network.VPN.private_key='4LmCvCIgkU1qDTuiNV2IRCwVr0/cH+fotYvuBxkqGVY='
		set network.VPN.listen_port='__WRTLAB_VPN_PORT__'
		add_list network.VPN.addresses='__WRTLAB_VPN_IP__/__WRTLAB_VPN_MASK_SHORT__'
	#VPN - Peers DEBUG ONLY
        	set network.@wireguard_VPN[-1].description='testclient'
        	set network.@wireguard_VPN[-1].public_key='wI/ESFQFyLyVGYTlT+mPAseC/Ly6rfDNgicRJkOCq2I='
        	add_list network.@wireguard_VPN[-1].allowed_ips='__WRTLAB_VPN_NET__'

#Network - Firewall - Traffic rules
        #Rule 'Allow-DHCP-renew"
	delete firewall.@rule[0]
        #Rule 'Allow-Wireguard'
	#add firewall rule
	#set firewall.@rule[-1].name='Allow-Wireguard'
	#set firewall.@rule[-1].family='ipv4'
	#set firewall.@rule[-1].src='wan'
	#set firewall.@rule[-1].dest_port='__WRTLAB_VPN_PORT__'
	#set firewall.@rule[-1].target='ACCEPT'
	#add_list firewall.@rule[-1].proto='udp'
        #Rule 'Allow-Iperf3'
	add firewall rule
	set firewall.@rule[-1].name='Allow-Iperf3'
	set firewall.@rule[-1].family='ipv4'
	set firewall.@rule[-1].src='*'
	set firewall.@rule[-1].dest_port='5201'
	set firewall.@rule[-1].target='ACCEPT'
        #Rule 'Allow-uhttpd' DEBUG ONLY
	add firewall rule
	set firewall.@rule[-1].name='Allow-Uhttpd'
	set firewall.@rule[-1].family='ipv4'
	set firewall.@rule[-1].src='*'
	set firewall.@rule[-1].dest_port='443'
	set firewall.@rule[-1].target='ACCEPT'
	add_list firewall.@rule[-1].proto='tcp'

#Network - Firewall - Zones - LAN
	delete firewall.@zone[0]

#Network - Firewall - Zones - VPN
	add firewall zone
	set firewall.@zone[-1].name='vpn'
	set firewall.@zone[-1].input='ACCEPT'
	set firewall.@zone[-1].output='ACCEPT'
	set firewall.@zone[-1].forward='REJECT'
	add_list firewall.@zone[-1].network='VPN'
	add firewall forwarding
	set firewall.@forwarding[-1].src='vpn'
	set firewall.@forwarding[-1].dest='wan'

#Network - Firewall - Port forwards
	#wireguard
	add firewall redirect
	set firewall.@redirect[-1].dest='wan'
	set firewall.@redirect[-1].target='DNAT'
	set firewall.@redirect[-1].name='wireguard'
	add_list firewall.@redirect[-1].proto='udp'
	set firewall.@redirect[-1].src='wan'
	set firewall.@redirect[-1].src_dport='__WRTLAB_VPN_PORT_HIDDEN__'
	set firewall.@redirect[-1].dest_ip='__WRTLAB_WAN_IP__'
	set firewall.@redirect[-1].dest_port='__WRTLAB_VPN_PORT__'
EOI
uci commit
